# hamclock/rootfs/etc/services.d/hamclock/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

## Get the 'message' key from the user config options.
message="message"

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message:="Starting HamClock..."}"

## Run your program
exec /usr/local/bin/hamclock -o
exec /hamclock/ESPHamClock/usr/local/bin/hamclock -o
## We don't know the correct path, yet, so for now run it like this

## After program started
bashio::log.info "${message:="HamClock started."}"
